# profile.d/k.sh
echo "Running KVM script"
/app/.k/kvm/kvm.sh
