# profile.d/k.sh
echo "Running KVM script"
sh /app/.k/kvm/kvm.sh
