# profile.d/libuv.sh
export LD_LIBRARY_PATH="/app/libuv/lib:${LD_LIBRARY_PATH}"