# profile.d/libuv.sh
echo "Adding libuv to LD_LIBRARY_PATH"
export LD_LIBRARY_PATH="/app/libuv/lib:${LD_LIBRARY_PATH}"