# profile.d/mono.sh
echo "Adding mono to PATH"
export PATH="/app/mono/bin:${PATH}"
export LD_LIBRARY_PATH="/app/mono/lib:${LD_LIBRARY_PATH}"