# profile.d/mono.sh
echo "Adding mono to PATH"
export PATH="/app/mono/bin:${PATH}"